# daemon/com.user.claude-usage-daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.claude-usage-daemon</string>

    <!-- The Swift daemon is a natively-signed Mach-O with an embedded Info.plist
         (NSBluetoothAlwaysUsageDescription + CFBundleIdentifier), so macOS
         attributes Bluetooth (TCC) access to it directly — no python/.app
         wrapper needed to own the permission prompt. -->
    <key>ProgramArguments</key>
    <array>
        <string>__DAEMON_BIN__</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__REPO_DIR__</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>__LOG_OUT__</string>

    <key>StandardErrorPath</key>
    <string>__LOG_ERR__</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>__HOME__</string>
        <key>PATH</key>
        <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
